# apps/tiny-tiny-rss/supervisord.conf
[program:tiny-tiny-rss]
command=bash -c "while true; do php update_daemon2.php; sleep 2; done"
directory=/data/http
user=www-data
