# apps/tiny-tiny-rss/supervisord.conf
[program:tiny-tiny-rss]
command=php update_daemon2.php
directory=/data/http
user=www-data
