# services/openvpn/supervisord.conf
[program:openvpn]
command=openvpn
directory=/data
