# frameworks/easyengine/supervisord.conf
[program:easyengine]
command=bash -c ""
