# frameworks/easyengine/supervisord.conf
[program:easyengine]
command=bash -c "ee stack start; tail -f /var/log/ee/ee.log"
