# services/subsonic/supervisord.conf
[program:subsonic]
command=java -Xmx256M -Djava.awt.headless=true -Dsubsonic.home=/data -jar subsonic-booter-jar-with-dependencies.jar
directory=/data
user=core
