# apps/adminer/supervisord.conf
[program:adminer]
command=php -S 0.0.0.0:8080
directory=/data
user=http
