# apps/adminer/supervisord.conf
[program:adminer]
command=/usr/bin/php -S 0.0.0.0:8080
directory=/config
user=http
