# applications/webirc/supervisord.conf
[program:webirc]
command=npm start
directory=/data
user=user
