# services/couchpotato/supervisord.conf
[program:couchpotato]
command=python CouchPotato.py --data_dir=/data
directory=/app
user=core
