# services/plexmediaserver/supervisord.conf
[program:avahi-daemon]
command=avahi-daemon

[program:dbus-daemon]
command=dbus-daemon --nofork --system

[program:plexmediaserver]
command=start_pms
environment="HOME=/data"
user=core
