# services/plexmediaserver/supervisord.conf
[program:avahi-daemon]
command=avahi-daemon

[program:plexmediaserver]
command=
