# services/sickbeard/supervisord.conf
[program:sickbeard]
command=python SickBeard.py --datadir=/data
directory=/app
user=core
