# service/nginx-php/supervisord.conf
[program:php-fpm]
command=php-fpm --fpm-config /etc/php5/fpm/php-fpm.conf

[program:php-fpm-config]
command=bash -c "while inotifywait -e create,delete,modify,move -q /data/config/php/; do supervisorctl restart php-fpm; done"
