# apps/haproxy/supervisord.conf
[program:haproxy]
command=haproxy -f /etc/haproxy/haproxy.cfg

[program:haproxy-config]
command=docker-gen -config /config/docker-gen.conf
