# apps/haproxy/supervisord.conf
[program:docker-gen]
command=docker-gen -config /etc/docker-gen.conf

[program:haproxy]
command=haproxy
