# apps/haproxy/supervisord.conf
[program:haproxy]
command=haproxy -f /data/haproxy.cfg

[program:haproxy-config]
command=bash -c "while inotifywait -e create,modify -q /data/haproxy.cfg; do supervisorctl restart haproxy; done"
