# apps/haproxy/supervisord.conf
[program:docker-gen]
command=/usr/local/bin/docker-gen -config /etc/docker-gen.conf

[program:haproxy]
command=/usr/bin/haproxy -sf $(cat /var/run/haproxy.pid)
