# services/haproxy/supervisord.conf
[program:haproxy]
command=haproxy -f /data/haproxy.cfg

[program:watch-config]
command=bash -c "while inotifywait -e create,modify,move -q /data/; do supervisorctl restart haproxy; done"
