# applications/dnsmasq/supervisord.conf
[program:dnsmasq]
command=dnsmasq

[program:dnsmasq-config]
command=bash -c "while inotifywait -e create,modify,move -q /data/dnsmasq/; do supervisorctl restart dnsmasq; done"
