# frameworks/nginx/supervisord.conf
[program:nginx]
command=nginx

[program:watch-config]
command=bash -c "while inotifywait -e create,delete,modify,move -q /config/; do supervisorctl restart nginx; done"
