# frameworks/nginx/supervisord.conf
[program:nginx]
command=nginx

[program:nginx-config]
command=bash -c "while inotifywait -e create,delete,modify,move -q /data/config/; do supervisorctl restart nginx; done"
