# service/postgresql/postgresql.conf
[program:postgresql]
command=/usr/lib/postgresql/9.3/bin/postgres -D /etc/postgresql/9.3/main
user=postgres
