# service/postgresql/supervisord.conf
[program:postgresql]
command=/usr/lib/postgresql/9.3/bin/postgres
user=postgres
