# services/haproxy-config/supervisord.conf
[program:docker-gen]
command=docker-gen -only-exposed -watch /config/haproxy.tmpl /data/haproxy.cfg
