# services/sniproxy/supervisord.conf
[program:sniproxy]
command=sniproxy -f

[program:sniproxy-config]
command=bash -c "while inotifywait -e create,modify,move -q /data/; do supervisorctl restart sniproxy; done"
