# services/sabnzbd/supervisord.conf
[program:sabnzbd]
command=sabnzbdplus -b 0 −f /data 
