# apps/znc/supervisord.conf
[program:znc]
command=znc -d /data -f
user=znc
