# apps/znc/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:znc]
command=/usr/bin/znc -d /data/znc -f
user=nobody
