# apps/znc/supervisord.conf
[program:znc]
command=/usr/bin/znc -d /data -f
user=znc
