# apps/znc/znc.conf
[program:znc]
command=/usr/bin/znc -d /data/znc -f
user=nobody
