# services/minecraft/supervisord.conf
[program:minecraft]
command=java -Xmx1024M -jar minecraft_server.jar
directory=/data
user=core
