# services/minecraft/supervisord.conf
[program:minecraft]
command=java -Xms1024M -Xmx1024M -jar minecraft_server.jar nogui
directory=/data
user=core
