# service/mariadb/mariadb.conf
[program:mariadb]
command=/usr/bin/mysqld_safe
