# apps/haproxy-config/supervisord.conf
[program:docker-gen]
command=docker-gen -notify "cp /config/data/haproxy.cfg /config/temp" -only-exposed -watch /config/data/haproxy.tmpl /config/data/haproxy.cfg
