# apps/haproxy-config/supervisord.conf
[program:docker-gen]
command=docker-gen -notify "cp /etc/haproxy.cfg /data/haproxy.cfg" -only-exposed -watch /etc/haproxy.tmpl /etc/haproxy.cfg
