# apps/haproxy-config/supervisord.conf
[program:docker-gen]
command=docker-gen -only-exposed -watch /config/data/haproxy.tmpl /config/temp/haproxy.cfg
