# apps/haproxy-config/supervisord.conf
[program:docker-gen]
command=docker-gen -config /config/docker-gen.conf
