# applications/tiny-tiny-rss/supervisord.conf
[program:tiny-tiny-rss]
command=php update.php --daemon
directory=/data/http
user=core
