# services/postgresql/supervisord.conf
[program:postgresql]
command=/usr/lib/postgresql/9.4/bin/postgres --config-file=/etc/postgresql/9.4/main/postgresql.conf
user=core
