# services/postgresql/supervisord.conf
[program:postgresql]
command=/usr/lib/postgresql/9.3/bin/postgres --config-file=/etc/postgresql/9.3/main/postgresql.conf
user=postgres
