# apps/minecraft/supervisord.conf
[program:minecraft]
command=/usr/bin/java -Xms512M -Xmx512M -jar minecraft_server.jar nogui
directory=/data
user=minecraft
